# rules/sigma/rules/Persistence/Linux/PER-E-1235-Detect Suspicious Modifications to Time Providers Registry Keys [high].yml
title: Detect Suspicious Modifications to Time Providers Registry Keys
status: Experimental
description: Detects modifications to registry keys used to configure Windows Time Providers.
author: Amr Ashraf
logsource:
  category: registry
  product: windows
detection:
  selection:
    EventID: 4657  # Registry modification event
    RegistryKey|contains: 'HKLM\\SYSTEM\\CurrentControlSet\\Services\\W32Time\\TimeProviders'
    RegistryValueName|contains:
      - 'DLLName'
      - 'Enabled'
  condition: selection
fields:
  - RegistryKey
  - RegistryValueName
  - NewValue
  - User
  - ComputerName
level: high
tags:
  - attack.persistence
  - attack.privilege_escalation
  - t1547.003
falsepositives:
  - Legitimate updates by administrators or system patches
