# Layers/EXE/EXE-004/EXE-004-E-000-0018.yml
title: WMI Creating Scheduled Tasks - EXE-004-E-000-0018
id: EXE-004-E-000-0018
status: experimental
description: |
  Detection Description:
  Detects scheduled task creation where the task is created via a command executed through WMIC process creation.

  Status:
  experimental

  Threat Mapping:
  MITRE ATT&CK Technique(s): T1053.005 - Scheduled Task/Job: Scheduled Task
  Kill Chain Phase: Execution

  Data Sources:
  MITRE Data Source: Process Creation
  Product targeted by the Rule: Windows
  Service in the product that provides the log: Sysmon (Event ID 1) or equivalent process creation logs

  Detection Logic:
  Matches WMIC process call create where the embedded command line invokes schtasks /create.

  False Positive Considerations:
  See falsepositives section below.

  Severity and Risk Rating:
  Severity Level: High

  References:
  See references section below.

  Author:
  HXP Detection team

  Use Case Sheet Details:
- L3 use cases per L2: 36
- L1 Use Case Name: Execution
- L1 Use Case Identifier: EXE
- L2 Use Case Name: WMI Execution
- L2 Use Case Identifier: EXE-004
- Rule Identifier: EXE-004-E-000-0018
- Technical Use Case Name: WMI Creating Scheduled Tasks
- Use Case Description: 
- Log source type: Windows Servers
- Detection Technology: SIEM
- Scope: All assets
author: HXP Detection team
date: 2025-12-13
references:
  - https://attack.mitre.org/techniques/T1053/005/
logsource:
  category: process_creation
  product: windows
  service: sysmon
detection:
  selection:
    Image|endswith: '\\wmic.exe'
    CommandLine|contains|all:
      - 'process'
      - 'call'
      - 'create'
    CommandLine|contains:
      - 'schtasks'
      - '/create'

  condition: selection

falsepositives:
  - Rare; could be admin automation chaining WMIC and schtasks
level: high
tags:
  - attack.persistence
  - attack.execution
  - attack.t1053.005
