# Threat_Profiling/Muddy Water Group/PER-6069-Anomalous Scheduled Tasks for Persistence [high].yml
title: Anomalous Scheduled Tasks for Persistence
status: test
description: Detects unusual scheduled task creation or modification that deviates from normal behavior, reducing false positives.
author: Amr Ashraf
references:
    - https://attack.mitre.org/techniques/T1053.005/
tags:
    - attack.persistence
    - attack.t1053.005
logsource:
    category: process_creation
    product: windows
detection:
    selection_schtasks:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains:
            - '/create'
            - '/change'
    selection_suspicious_paths:
        CommandLine|contains:
            - 'AppData\Roaming\'
            - 'Temp\'
            - 'C:\Users\Public\'
            - 'C:\Windows\Temp\'
            - 'C:\Windows\System32\Tasks\'
    selection_hidden_task:
        CommandLine|contains:
            - '/ru SYSTEM'
            - '/tn "svchost"'
            - '/tn "UpdateService"'
            - '/tn "WindowsSecurity"'
            - '/tn "DefenderUpdate"'
            - '/tn "WinLogonHelper"'
    selection_executable_in_task:
        CommandLine|contains:
            - 'cmd.exe /c'
            - 'powershell.exe -w hidden'
            - 'wscript.exe'
            - 'cscript.exe'
            - 'mshta.exe'
            - 'rundll32.exe'
            - 'regsvr32.exe'
    condition: selection_schtasks and (selection_suspicious_paths or selection_hidden_task or selection_executable_in_task)
falsepositives:
    - Legitimate scheduled task modifications by administrators, security tools, or updates.
level: high
