# EXE/EXE-003/EXE-003-E-000-0003.yml
title: Cron Job Executing from Suspicious Directories - EXE-003-E-000-0003
id: EXE-003-E-000-0003
status: experimental
description: |
  Detection Description:
  Detects cron spawning processes from suspicious directories that are frequently used for staging payloads.

  Status:
  experimental

  Threat Mapping:
  MITRE ATT&CK Technique(s): T1053.003 - Scheduled Task/Job: Cron
  Kill Chain Phase: Execution

  Data Sources:
  MITRE Data Source: Process Creation
  Product targeted by the Rule: Linux
  Service in the product that provides the log: auditd / syslog process execution (as ingested to SIEM)

  Detection Logic:
  Matches Linux process creation events where cron/crond is the parent and the executed file path is in /tmp, /dev/shm, /var/tmp, or other user-writable locations.

  False Positive Considerations:
  See falsepositives section below.

  Severity and Risk Rating:
  Severity Level: High

  References:
  See references section below.

  Author:
  HXP Detection team

  Use Case Sheet Details:
- L3 use cases per L2: 4
- L1 Use Case Name: Execution
- L1 Use Case Identifier: EXE
- L2 Use Case Name: Scheduled Task Execution
- L2 Use Case Identifier: EXE-003
- Rule Identifier: EXE-003-E-000-0003
- Technical Use Case Name: Cron Job Executing from Suspicious Directories
- Use Case Description: 
- Log source type: Linux Servers
- Detection Technology: SIEM
- Scope: All assets
author: HXP Detection team
date: 2025-12-13
references:
  - https://attack.mitre.org/techniques/T1053/003/
logsource:
  category: process_creation
  product: linux
  service: auditd
detection:
  selection:
    ParentImage|endswith:
      - '/cron'
      - '/crond'
    Image|contains:
      - '/tmp/'
      - '/var/tmp/'
      - '/dev/shm/'
      - '/run/user/'
      - '/home/'
  condition: selection

falsepositives:
  - Legitimate admin scripts executed from /tmp during maintenance (should be rare; tune with allowlists)
level: high
tags:
  - attack.execution
  - attack.t1053.003
