# PER/PER-003/PER-003-E-000-0002.yml
title: Scheduled Task Modified - PER-003-E-000-0002
id: PER-003-E-000-0002
status: experimental
description: |
      Detection Description:
      Detects modification of existing scheduled tasks, which can indicate persistence or defense evasion.

      Status:
      experimental

      Threat Mapping:
      MITRE ATT&CK Technique(s): T1053.005 - Scheduled Task/Job: Scheduled Task
      Kill Chain Phase: Persistence

      Data Sources:
      MITRE Data Source: Scheduled Task
      Product targeted by the Rule: Windows
      Service in the product that provides the log: Security (4702) and/or TaskScheduler Operational logs

      Detection Logic:
      Windows Security 4702 indicates a scheduled task was updated.

      False Positive Considerations:
      Legitimate task updates during patching or application upgrades.

      Severity and Risk Rating:
      Severity Level: Medium

      References:
      See references section below.

      Author:
      HXP Detection team

      Use Case Sheet Details:
- L3 use cases per L2: 9
- L1 Use Case Name: Persistence
- L1 Use Case Identifier: PER
- L2 Use Case Name: Scheduled Task Persistence
- L2 Use Case Identifier: PER-003
- Rule Identifier: PER-003-E-000-0002
- Technical Use Case Name: Scheduled Task Modified
- Use Case Description: 
- Log source type: Windows Servers
- Detection Technology: SIEM
- Scope: All assets
author: HXP Detection team
date: 2025-12-13
references:
  - https://attack.mitre.org/techniques/T1053/005/
logsource:
  product: windows
  service: security
detection:
  selection:
    EventID: 4702
  condition: selection
falsepositives:
  - Legitimate scheduled task updates
level: medium
tags:
  - attack.persistence
  - attack.t1053.005
