# PER/PER-003/PER-003-E-000-0007.yml
title: AT.exe Command Execution for Job Scheduling - PER-003-E-000-0007
id: PER-003-E-000-0007
status: experimental
description: |
      Detection Description:
      Detects use of at.exe to schedule jobs. AT is deprecated but still abused for persistence/execution.

      Status:
      experimental

      Threat Mapping:
      MITRE ATT&CK Technique(s): T1053.002 - Scheduled Task/Job: At
      Kill Chain Phase: Persistence

      Data Sources:
      MITRE Data Source: Process Creation
      Product targeted by the Rule: Windows
      Service in the product that provides the log: Sysmon (Event ID 1) or equivalent

      Detection Logic:
      at.exe process creation.

      False Positive Considerations:
      Rare legacy administration; should be uncommon in modern environments.

      Severity and Risk Rating:
      Severity Level: High

      References:
      See references section below.

      Author:
      HXP Detection team

      Use Case Sheet Details:
- L3 use cases per L2: 9
- L1 Use Case Name: Persistence
- L1 Use Case Identifier: PER
- L2 Use Case Name: Scheduled Task Persistence
- L2 Use Case Identifier: PER-003
- Rule Identifier: PER-003-E-000-0007
- Technical Use Case Name: AT.exe Command Execution for Job Scheduling
- Use Case Description: 
- Log source type: Windows Servers
- Detection Technology: SIEM
- Scope: All assets
author: HXP Detection team
date: 2025-12-13
references:
  - https://attack.mitre.org/techniques/T1053/002/
logsource:
  category: process_creation
  product: windows
  service: sysmon
detection:
  selection:
    Image|endswith: '\at.exe'
  condition: selection
falsepositives:
  - Legacy administrative scheduling in older environments
level: high
tags:
  - attack.persistence
  - attack.t1053.002
