# Layers/PER/PER-003/PER-003-E-000-0009.yml
title: AT Remote Job Scheduling - PER-003-E-000-0009
id: PER-003-E-000-0009
status: experimental
description: |
      Detection Description:
      Detects remote job scheduling with at.exe targeting a remote host (\\HOST), often used for lateral movement and persistence.

      Status:
      experimental

      Threat Mapping:
      MITRE ATT&CK Technique(s): T1053.002 - At, T1021 - Remote Services
      Kill Chain Phase: Persistence

      Data Sources:
      MITRE Data Source: Process Creation
      Product targeted by the Rule: Windows
      Service in the product that provides the log: Sysmon (Event ID 1) or equivalent

      Detection Logic:
      at.exe command line contains a remote host specification (\\).

      False Positive Considerations:
      Rare; legitimate remote administration on legacy systems.

      Severity and Risk Rating:
      Severity Level: High

      References:
      See references section below.

      Author:
      HXP Detection team

      Use Case Sheet Details:
- L3 use cases per L2: 9
- L1 Use Case Name: Persistence
- L1 Use Case Identifier: PER
- L2 Use Case Name: Scheduled Task Persistence
- L2 Use Case Identifier: PER-003
- Rule Identifier: PER-003-E-000-0009
- Technical Use Case Name: AT Remote Job Scheduling
- Use Case Description: 
- Log source type: Windows Servers
- Detection Technology: SIEM
- Scope: All assets
author: HXP Detection team
date: 2025-12-13
references:
  - https://attack.mitre.org/techniques/T1053/002/
logsource:
  category: process_creation
  product: windows
  service: sysmon
detection:
  selection:
    Image|endswith: '\at.exe'
    CommandLine|contains: '\\\\'
  condition: selection
falsepositives:
  - Legacy remote administration
level: high
tags:
  - attack.persistence
  - attack.lateral_movement
  - attack.t1053.002
